# data/units/enlightenment.service.in
[Unit]
Description=Enlightenment service
Before=end.target
After=xorg.target
Requires=xorg.target
Requires=dbus.socket
Wants=pulseaudio.service
Wants=ssh-agent.service
AllowIsolate=true

[Service]
Type=notify
#Environment=PATH=uncomment:to:override:your:PATH
Environment=E_START=enlightenment
ExecStart=@prefix@/bin/enlightenment
Restart=always
RestartPreventExitStatus=0
RestartSec=2
StartLimitInterval=30
StartLimitBurst=7
KillMode=process
WatchdogSec=10
TimeoutSec=5
NotifyAccess=all

[Install]
WantedBy=wm.target
